# _config/autostart/compton.desktop
[Desktop Entry]
Type=Application
Name=Compton
Exec=compton
X-GNOME-Autostart-enabled=true
Comment=Composition manager
