# _config/autostart/kalu.desktop
[Desktop Entry]
Type=Application
Name=Kalu
Exec=kalu
X-GNOME-Autostart-enabled=true
Comment=Kalu
