# _config/autostart/keepassxc.desktop
[Desktop Entry]
Type=Application
Name=KeePassXC
Exec=keepassxc
X-GNOME-Autostart-enabled=true
Comment=KeePassXC
