# _config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network Manager
Exec=nm-applet
X-GNOME-Autostart-enabled=true
Comment=Network Manager
