# _config/autostart/syncthing.desktop
[Desktop Entry]
Type=Application
Name=Syncthing
Exec=syncthing-gtk -m
X-GNOME-Autostart-enabled=true
Comment=Syncthing
