# _config/autostart/xmodmap.desktop
[Desktop Entry]
Type=Application
Name=xmodmap
Exec=xmodmap .xmodmap
X-GNOME-Autostart-enabled=true
Comment=Load xmodmap
