# _config/autostart/xpad.desktop
[Desktop Entry]
Type=Application
Name=zim
Exec=zim
X-GNOME-Autostart-enabled=true
Comment=Start zim
