# _config/autostart/yubioath.desktop
[Desktop Entry]
Type=Application
Name=yubioath-desktop
Exec=yubioath-desktop
X-GNOME-Autostart-enabled=true
Comment=Yubioath
